# templates/runit/log/run
#!/bin/sh

exec chpst -u{{ services_logger_user }} svlogd -tt ./main
