# templates/runit/log/run
#!/bin/sh

exec chpst -u{{ service_logger_user }} svlogd -tt ./main
